# deploy/sabji-gpt.service
[Unit]
Description=SabjiGPT MCP Server and Scheduler
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/sabjiGPT
Environment=PATH=/opt/sabjiGPT/venv/bin
Environment=AUTH_TOKEN=sabji_gpt_secret_2025
Environment=MCP_PORT=8087
Environment=SCRAPE_HEADLESS=true
Environment=LOG_LEVEL=INFO
Environment=DATABASE_PATH=/opt/sabjiGPT/data/mandi_prices.db
ExecStart=/opt/sabjiGPT/venv/bin/python run_automated_system.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
